# linux_config/chennode.service
[Unit]
Description=chen-node booter

[Service]
WorkindDirectory=/root/ChenNode/
ExecStart=/usr/local/bin/forever index.js

